# puphpet/puppet/modules/rabbitmq/templates/rabbitmq-server.service.d/limits.conf
[Service]
LimitNOFILE=<%= @file_limit %>
